# systemd/church-bells@.service
[Unit]
Description=Church Bell System - Speaker Control Service
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/church-music-app
Environment=PATH=/home/%i/church-music-app/venv/bin:/usr/local/bin:/usr/bin:/bin

# PipeWire / WirePlumber run as user-session services.
# VLC and wpctl need XDG_RUNTIME_DIR and DBUS_SESSION_BUS_ADDRESS to reach
# the PipeWire daemon.  We use a shell wrapper to resolve the UID at runtime
# because %U is the manager UID (root), not the User= UID.
ExecStart=/bin/bash -c '\
  export XDG_RUNTIME_DIR="/run/user/$$(id -u)" && \
  export DBUS_SESSION_BUS_ADDRESS="unix:path=/run/user/$$(id -u)/bus" && \
  export PULSE_SERVER="unix:/run/user/$$(id -u)/pulse/native" && \
  exec /home/%i/church-music-app/venv/bin/python run.py'

Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=church-bells

# Give access to audio and bluetooth
SupplementaryGroups=audio bluetooth

[Install]
WantedBy=multi-user.target
